# astrolog-api.service
[Unit]
Description=Astrolog API Server
After=network.target caddy.service
Wants=caddy.service

[Service]
Type=simple
User=ruslan
WorkingDirectory=/home/ruslan/aia/server
# Load environment variables from .env file (if exists)
EnvironmentFile=-/home/ruslan/aia/server/.env
# Override PORT to 8080 (Caddy will handle HTTPS on 443)
Environment="PORT=8080"
ExecStart=/home/ruslan/aia/server/astrolog_api
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
