# astrolog-api.service
[Unit]
Description=Astrolog API Server with HTTPS
After=network.target

[Service]
Type=simple
User=ruslan
WorkingDirectory=/home/ruslan/aia/server
EnvironmentFile=/home/ruslan/aia/server/.env
ExecStart=/home/ruslan/aia/server/astrolog_server
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
